# server/deploy/services/act-keys-import/run
#! /bin/bash

(cd $ACT_APP_DIR && bundle exec rake act:keys_import)